# bin/host-bridge/com.missioncontrol.terminalbridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent TEMPLATE for uninterrupted bridge runtime (survives reboot).
  Do NOT load directly: install-agent.sh (v16) renders the templated values
  (repo path, bridge.token, node path via `command -v node`), writes 0600 to
  ~/Library/LaunchAgents/, and runs launchctl bootstrap gui/$UID + kickstart.
  Auto-invoked by install.sh; removed by uninstall.sh. Manual: make install-agent
  (DRY_RUN=1 make install-agent to rehearse).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.missioncontrol.terminalbridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>REPLACE_WITH_ABSOLUTE_PATH/bin/host-bridge/terminal-bridge.mjs</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>7345</string>
    <key>HOST_BRIDGE_TOKEN</key>
    <string>REPLACE_WITH_TOKEN_FROM_~/.llm-mission-control/bridge.token</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <!-- Self-healing: restart only on crash (non-zero exit), throttle restart
       storms to >=10s apart. Pattern: tjluoma/launchd-keepalive (public domain). -->
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/tmp/llm-bridge/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/llm-bridge/launchd.err.log</string>
</dict>
</plist>
