# scripts/com.ollamas.brain-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.ollamas.brain-loop</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string><string>-lc</string>
    <string>cd $HOME/Desktop/ollamas &amp;&amp; npx tsx scripts/brain-loop.ts</string>
  </array>
  <key>StartInterval</key><integer>900</integer>
  <key>RunAtLoad</key><false/>
  <key>KeepAlive</key><false/>
  <key>Nice</key><integer>10</integer>
  <key>StandardOutPath</key><string>/tmp/ollamas-brain-loop.log</string>
  <key>StandardErrorPath</key><string>/tmp/ollamas-brain-loop.err</string>
</dict>
</plist>
