# scripts/com.ollamas.brain-obsidian-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Obsidian mirror sync: every 5 min reconcile brain ⇄ ~/ollamas-vault via the :3000
     HTTP surface (thin client; never opens brain.db). Idempotent + cheap (no-change tick
     writes 0). Mirrors com.ollamas.brain-loop (proven zsh form on this machine).
     Activate once:  cp scripts/com.ollamas.brain-obsidian-sync.plist ~/Library/LaunchAgents/ &&
                     launchctl load -w ~/Library/LaunchAgents/com.ollamas.brain-obsidian-sync.plist -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.ollamas.brain-obsidian-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string><string>-lc</string>
    <string>cd $HOME/Desktop/ollamas &amp;&amp; npx tsx scripts/brain-obsidian-sync.ts</string>
  </array>
  <key>StartInterval</key><integer>300</integer>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><false/>
  <key>Nice</key><integer>10</integer>
  <key>StandardOutPath</key><string>/tmp/ollamas-obsidian-sync.log</string>
  <key>StandardErrorPath</key><string>/tmp/ollamas-obsidian-sync.err</string>
</dict>
</plist>
